# buildroot-external/rootfs-overlay/etc/systemd/system/NetworkManager-wait-online.service.d/distro.conf
[Service]
ExecStart=/usr/bin/nm-online -q --timeout=3600
