# buildroot-external/rootfs-overlay/etc/systemd/system/smb.service.d/distro.conf
[Unit]
RequiresMountsFor=/var/lib/samba
ConditionPathExists=/var/lib/samba/enabled
