# buildroot-external/rootfs-overlay/etc/systemd/system/wsdd2.service.d/distro.conf
[Service]
DynamicUser=false
ProtectSystem=full