# buildroot-external/rootfs-overlay/etc/systemd/system/wsdd2.service.d/distro.conf
[Service]
AmbientCapabilities=CAP_NET_RAW
ProtectSystem=full