# buildroot-external/rootfs-overlay/usr/lib/systemd/system/channels-dvr.service
[Unit]
Description=Channels DVR
Wants=network-online.target time-sync.target
After=network-online.target time-sync.target
RequiresMountsFor=/mnt/data /media

[Service]
Type=simple
User=channels
Group=channels
WorkingDirectory=/mnt/data/channels-dvr/data
Environment=CHANNELS_DVR_DEFAULT_PATH=/media/DVR
ExecStart=/mnt/data/channels-dvr/latest/channels-dvr -dir /mnt/data/channels-dvr/data/
Restart=always
RestartSec=30
LimitNOFILE=8192
OOMScoreAdjust=-500

[Install]
WantedBy=multi-user.target
