# buildroot-external/rootfs-overlay/etc/systemd/system/dropbear.service.d/distro.conf
[Install]
WantedBy=multi-user.target
WantedBy=recovery.target
